# docker/cron/crontab/docker



* * * * * /usr/local/bin/php /var/www/src/webCrawler/artisan pop-url >> /var/log/cron.log
# * * * * * /bin/echo "foobar" >> /var/log/cron.log
